# .deploy/madentk-chat-agent-ui.service
[Unit]
Description=Chat Agent Streamlit UI Service
After=network.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__APP_DIR__
Environment=PATH=__APP_DIR__/madentk-chat-agent-venv/bin
EnvironmentFile=/etc/madentk-chat-agent-ui.env
ExecStart=/bin/bash __APP_DIR__/start-ui.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=madentk-chat-agent-ui

[Install]
WantedBy=multi-user.target
