# .deploy/madentk-chat-agent.service
[Unit]
Description=Chat Agent FastAPI Service
After=network.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__APP_DIR__
Environment=PATH=__APP_DIR__/madentk-chat-agent-venv/bin
EnvironmentFile=/etc/madentk-chat-agent.env
ExecStart=__APP_DIR__/madentk-chat-agent-venv/bin/uvicorn app:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=madentk-chat-agent

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__APP_DIR__ /var/log/madentk-chat-agent /data

[Install]
WantedBy=multi-user.target
